# rc.local
chmod +x /BashBunny/HID.sh
chmod +x /BashBunny/apache.sh
chmod +x /BashBunny/UPDATE_WEBSERVER.sh
chmod +x /BashBunny/ETHERNET.sh


exit 0
